# systemd/codex-telegram-bridge.service.template
[Unit]
Description=Codex Telegram Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__REPO_ROOT__
Environment=PYTHONUNBUFFERED=1
Environment=CODEX_TELEGRAM_BRIDGE_CONFIG=__CONFIG_PATH__
ExecStart=__PYTHON_BIN__ __REPO_ROOT__/telegram_codex_relay/telegram_bot.py
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
